# letsencrypt_container_service.service
[Unit]
Description=Update SSL certificates with Certbot

[Service]
Type=oneshot
ExecStart=/usr/local/bin/update_certs.sh
